# server_setup/caddy.service
[Unit]
Description=Caddy - Treekipedia GraphFlow HTTPS Proxy
Documentation=https://caddyserver.com/docs/
After=network.target network-online.target
Requires=network-online.target

[Service]
Type=notify
User=root
Group=root
ExecStart=/usr/bin/caddy run --config /etc/caddy/caddy.json --adapter json
ExecReload=/usr/bin/caddy reload --config /etc/caddy/caddy.json --adapter json --force
TimeoutStopSec=5s
LimitNOFILE=1048576
LimitNPROC=512
PrivateTmp=true
ProtectSystem=full
AmbientCapabilities=CAP_NET_BIND_SERVICE
Environment="CADDY_ADMIN=0.0.0.0:2019"
StandardOutput=file:/var/log/caddy/caddy.log
StandardError=file:/var/log/caddy/caddy-error.log

[Install]
WantedBy=multi-user.target
